# files/forgejo-runner/forgejo-runner.service
[Unit]
Description=Forgejo Actions runner daemon (UM790)
Documentation=https://forgejo.tail5bbc4.ts.net
# Image-baked default: NOT enabled, NOT started.
# Per the user's spec ("install + configure, not start, not
# authorize"). The operator runs /usr/local/bin/setup-forgejo-runner.sh
# which:
#   1. Verifies rootless-docker is up + groups render/video present.
#   2. Materializes the per-user token + config files.
#   3. Runs `systemctl --user enable --now forgejo-runner.service`.
#
# Until then, this unit is dormant. Systemd without [Install] = nothing
# auto-starts it. `systemctl --user enable forgejo-runner.service` adds
# a symlink in the user's wanted-by directory; this is what flips it on.

[Service]
Type=simple
# User-mode systemd units implicitly run as the calling user. No
# User=/Group= required. (Contrast with system-mode units, which would
# need them.)
WorkingDirectory=%h/.local/share/forgejo-runner

# Wait for the rootless docker socket before starting the daemon.
# Reuses the same shape as setup-rootless-docker's verify script.
# `--wait-docker` mode is a no-op except for polling — idempotent.
ExecStartPre=/usr/local/bin/setup-forgejo-runner.sh --wait-docker

ExecStart=/usr/bin/forgejo-runner daemon --config %h/.config/forgejo-runner/config.yaml

# DOCKER_HOST must be exported BEFORE the daemon reads config.yaml.
# The setup script writes /run/user/<UID>/forgejo-runner.env at
# runtime with the correct DOCKER_HOST line. systemd reads it here.
EnvironmentFile=-%h/.config/forgejo-runner/env

Restart=on-failure
RestartSec=20s
TimeoutStartSec=180

# Hardening (subset — rootless docker is already hardened at its
# daemon). ProtectSystem=full keeps the global FS read-only except
# /dev /proc /sys; we let the daemon write to /var/tmp/forgejo-runner
# via RunTimeDirectory-style via ProtectSystem=full alone; for v1 we
# accept that /etc is read-only and the runner can never write there.
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
# AF_UNIX needed for the docker socket; AF_INET* for Forgejo API +
# artifact transfer; AF_NETLINK needed by Docker for namespace setup.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
NoNewPrivileges=true
UMask=0077
